# systemd/s3-sync-hls.service
[Unit]
Description=Continuously sync HLS segments to S3
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
# Load env with BUCKET, PREFIX, SYNC_INTERVAL_SECONDS
EnvironmentFile=/etc/hls-to-s3.env
ExecStart=/usr/local/bin/s3-sync-hls.sh
Restart=always
RestartSec=2
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
